# etc/systemd/system/personal-drafts-sync.service
[Unit]
Description=Personal Drafts Sync Service
After=network.target

[Service]
Type=oneshot
User=%i
Environment=PD_DIR=%h/notes
ExecStart=/usr/local/bin/pd-pull
ExecStart=/usr/local/bin/pd-push

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
